# services/ai-player-bot.service
[Unit]
Description=AI Player Bot — onemoment-gateway TOS-compliant AI gateway
After=network.target

[Service]
Type=simple
User=ai-host-user
WorkingDirectory=/opt/onemoment-gateway
ExecStart=/usr/bin/python3 /opt/onemoment-gateway/bots/ai-player-bot.py
Restart=on-failure
RestartSec=10
EnvironmentFile=/opt/onemoment-gateway/.env

[Install]
WantedBy=multi-user.target
